# dev-environment/rest.crontab
# ┌───────── minute (0 - 59)
# │ ┌─────── hour (0 - 23)
# │ │ ┌───── day of the month (1 - 31)
# │ │ │ ┌─── month (1 - 12)
# │ │ │ │ ┌─ day of the week (0 - 6) (Sunday to Saturday)
# │ │ │ │ │
# * * * * * command to execute

  * * * * * python manage.py send_eboks

# For at teste logrotation. Den vil rotere logs hvert minut,
# så den er udkommenteret for at undgå at spamme filer i log-mappen
# * * * * * logrotate --force /logrotate.conf
